# Soal-3/crontab
5 6-23/8 * * 1-5,7 /bin/bash /home/evelyn/Praktikum1/soal3a.sh
